# observer/egenta.service
[Unit]
Description=Egenta, always-on observer of Claude Code sessions
Documentation=file:///home/keletonik/Egenta/README.md
After=default.target

[Service]
Type=simple
# Single instance; the entrypoint takes its own lock as a second guard.
# Stdlib only, so plain python3 needs no uv and no virtualenv. --daemon is the
# resident sweep loop; the deep brain stays dormant until ANTHROPIC_API_KEY is set.
ExecStart=/usr/bin/python3 %h/Egenta/observer/egenta.py --daemon --poll 120
WorkingDirectory=%h/Egenta
Restart=on-failure
RestartSec=5
# Resident but cheap. Cap it so a runaway analysis cannot eat the box.
MemoryMax=2G
CPUWeight=20
# Hardening: home is read-only except Egenta's own tree and its XDG state dir
# (the cursor and lock). It reads the session logs, writes nowhere else. The
# write gate, not the kernel, is what keeps the daemon out of ~/.claude.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/Egenta %h/.local/state

[Install]
WantedBy=default.target

# Install when the entrypoint exists:
#   ln -s ~/Egenta/observer/egenta.service ~/.config/systemd/user/egenta.service
#   systemctl --user daemon-reload
#   systemctl --user enable --now egenta.service
# Uninstall:
#   systemctl --user disable --now egenta.service
#   rm ~/.config/systemd/user/egenta.service && systemctl --user daemon-reload
